# docs/systemd/npm-proxy-guard.service.d/upstream.conf
# -*- coding: utf-8, tab-width: 2 -*-
#
# Explicitly set proxy, in case you don't have it in global env vars.

[Service]
Environment = TGT_HOST=localhost TGT_PORT=8118
